# misc/udev_rules/95-microsoft-natural-keyboard.rules
ENV{ID_VENDOR}=="Microsoft", ENV{ID_MODEL_ID}=="00db", RUN+="keymap $name 0xc022d up 0xc022e down"

